# init.d/xray
#!/sbin/openrc-run

name="Xray"
description="The best v2ray-core, with XTLS support"
description_checkconfig="Test configuration file"

command="/usr/local/bin/xray"
command_args="run -c /usr/local/etc/xray/config.json"
command_user="root"

pidfile="/run/xray.pid"
command_background="yes"

depend() {
	need net
}
